# templates/etc/systemd/system/systemd-journald.service.d/after-zfs-mount.service.j2
{{ ansible_managed | comment }}

[Unit]
After=zfs-mount.service
